# rhelai-validation/templates/ilab-serve.service.j2
[Unit]
Description=ilab model serve service

[Install]
WantedBy=multi-user.target default.target

[Service]
ExecStart=ilab model serve --model-path {{model_path | quote}}
Restart=always
